# sync/etc/systemd/system/puma.socket.erb
[Unit]
Description=Puma HTTP Server Accept Sockets

[Socket]
<% listen_streams.each do |s| %>
ListenStream=<%= s %>
<% end %>

# Socket options matching what Puma wants
NoDelay=true
ReusePort=true

[Install]
WantedBy=sockets.target
