# sync/etc/systemd/system/puma.socket.erb
[Unit]
Description=Puma HTTP Server Accept Sockets

[Socket]
<% listen_streams.each do |s| %>
ListenStream=<%= s %>
<% end %>

# Socket options matching Puma defaults
NoDelay=true
ReusePort=true
Backlog=1024

[Install]
WantedBy=sockets.target
