# sync/etc/systemd/system/puma.service.erb
[Unit]
Description=Puma HTTP Server
After=network.target
<% if systemd_socket %>
Requires=<%= systemd_socket %>
<% end %>

[Service]
Type=simple
User=<%= run_user %>
WorkingDirectory=<%= rack_path %>
ExecStart=<%= puma_start_command %>
Restart=always

[Install]
WantedBy=multi-user.target
